# br2-external/rootfs_overlay/etc/systemd/system/serial-getty@.service.d/reset-tty.conf
################################################################################
#         ____  ___________               __          __                       #
#        / __ \/ __/ __/ (_)___  ___     / /   ____ _/ /_                      #
#       / / / / /_/ /_/ / / __ \/ _ \   / /   / __ `/ __ \                     #
#      / /_/ / __/ __/ / / / / /  __/  / /___/ /_/ / /_/ /                     #
#      \____/_/ /_/ /_/_/_/ /_/\___/  /_____/\__,_/_.___/                      #
#                                                                              #
#      Copyright (C) 2025-2026 Offline Lab                                     #
#      Contact: info@offline-lab.com                                           #
#      SPDX-License-Identifier: AGPL-3.0-only                                  #
################################################################################

[Service]
# Reset terminal on restart so echo/icanon are sane after logout.
# The upstream unit uses --noreset which inherits the shell's raw-mode
# settings when agetty restarts; this drop-in removes that flag.
ExecStart=
ExecStart=-/sbin/agetty --noclear --issue-file=/etc/issue:/etc/issue.d:/run/issue.d:/usr/lib/issue.d --keep-baud 115200,57600,38400,9600 - $TERM
