# br2-external/package/offlinelab-usb-gadget/src/systemd/service/usb-gadget.service
[Unit]
Description=Setup USB composite gadget (ACM + ECM)
After=sys-kernel-config.mount
Requires=sys-kernel-config.mount

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/usb-gadget.sh

[Install]
WantedBy=sysinit.target
