# overlays/configure-system/etc/systemd/system/bluetooth.service.d/override.conf
[Service]
ConfigurationDirectoryMode=0755
