# br2-external/package/offlinelab-wifi/src/systemd/service/provision-wifi.service
[Unit]
Description=Provision WiFi config from boot partition
After=boot-firmware.mount expand-data.service
Wants=boot-firmware.mount

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/provision-wifi.sh

[Install]
WantedBy=multi-user.target
