# br2-external/package/offlinelab-resources/src/systemd/service/offlinelab-resources.service
[Unit]
Description=Offline Lab Resource Baseline Measurement
After=multi-user.target
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/init-resources

[Install]
WantedBy=multi-user.target
