# br2-external/package/offlinelab-portable/src/hello-world/usr/lib/systemd/system/hello-portable.service
[Unit]
Description=Hello portable test service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo "hello-portable: $(date)" > /run/hello-portable.stamp'
RemainAfterExit=yes
