# br2-external/package/offlinelab-disco/src/systemd/service/disco-daemon.service
################################################################################
#         ____  ___________               __          __                       #
#        / __ \/ __/ __/ (_)___  ___     / /   ____ _/ /_                      #
#       / / / / /_/ /_/ / / __ \/ _ \   / /   / __ `/ __ \                     #
#      / /_/ / __/ __/ / / / / /  __/  / /___/ /_/ / /_/ /                     #
#      \____/_/ /_/ /_/_/_/ /_/\___/  /_____/\__,_/_.___/                      #
#                                                                              #
#      Copyright (C) 2025-2026 Offline Lab                                     #
#      Contact: info@offline-lab.com                                           #
#      SPDX-License-Identifier: AGPL-3.0-only                                  #
################################################################################

[Unit]
Description=Disco - Service Discovery Daemon for offline networks
Documentation=https://github.com/offline-lab/disco
After=offlinelab-sysusers.service systemd-networkd.service
ConditionPathExists=/data/config/services/disco

[Service]
Type=notify
User=disco
Group=disco
ExecStart=/usr/bin/disco-daemon --config /data/config/disco/config.yaml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
RuntimeDirectory=disco
RuntimeDirectoryMode=0750
ConfigurationDirectory=disco
ConfigurationDirectoryMode=0750
StateDirectory=disco
StateDirectoryMode=0750
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SYS_TIME
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SYS_TIME
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/data/config/disco
ReadWritePaths=/run/disco /var/lib/disco
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
