# br2-external/package/offlinelab-bootconf/src/systemd/service/bootconf.service
################################################################################
#         ____  ___________               __          __                       #
#        / __ \/ __/ __/ (_)___  ___     / /   ____ _/ /_                      #
#       / / / / /_/ /_/ / / __ \/ _ \   / /   / __ `/ __ \                     #
#      / /_/ / __/ __/ / / / / /  __/  / /___/ /_/ / /_/ /                     #
#      \____/_/ /_/ /_/_/_/ /_/\___/  /_____/\__,_/_.___/                      #
#                                                                              #
#      Copyright (C) 2025-2026 Offline Lab                                     #
#      Contact: info@offline-lab.com                                           #
#      SPDX-License-Identifier: AGPL-3.0-only                                  #
################################################################################

[Unit]
Description=Bootconf — boot-time configuration
Documentation=https://github.com/offline-lab/bootconf
After=local-fs.target boot-firmware.mount
Wants=boot-firmware.mount
ConditionPathExists=/data/config/bootconf.yaml
Before=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/bootconf run --config /data/config/bootconf.yaml
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
