# br2-external/package/offlinelab-wifi/src/udev/80-wifi-power-save.rules
# Disable WiFi power save on wlan0 — prevents brcmfmac from going unresponsive
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/usr/sbin/iw dev wlan0 set power_save off"
